# files/etc/uci-defaults/999-fix-opkg-distfeeds
# Remove feeds that are not auto-built
sed -i '/ nss_packages /d' /etc/opkg/distfeeds.conf
sed -i '/ sqm_scripts_nss /d' /etc/opkg/distfeeds.conf